# ops/macos/com.revenuepilot.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.revenuepilot.backend</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <!-- Replace __REPO__ with your absolute repo path; venv preferred but poetry or system python also okay -->
    <string>cd "__REPO__" && .venv/bin/python -m uvicorn backend.main:app --host 127.0.0.1 --port 8000 --reload</string>
  </array>
  <key>WorkingDirectory</key><string>__REPO__</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>$HOME/Library/Logs/RevenuePilot/backend.log</string>
  <key>StandardErrorPath</key><string>$HOME/Library/Logs/RevenuePilot/backend.err</string>
</dict></plist>
